# supervisor/couchdb.conf
[program:couchdb]
command = /usr/bin/couchdb
priority = 1
autostart = true
autorestart = true
user = couchdb
