# supervisor/foreman.conf
[program:foreman]
command=/usr/local/bin/foreman start -d /app
